# nct6775_fancontrol/fancontrol.service
[Unit]
Description=nct6775_fancontrol

[Service]
ExecStart=/usr/bin/nct6775_fancontrol

[Install]
Alias=nct6775_fancontrol.service

